# atop-rotate.timer
[Unit]
Description=Daily atop restart

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
