# atop.cron
# start atop daily at midnight
0 0 * * * root /etc/atop/atop.daily
