# atop-rotate.service
[Unit]
Description=Restart atop daemon to rotate logs

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart atop.service
