# atop-rotate.timer
[Unit]
Description=Daily atop restart
Documentation=man:atop(1)

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
